# http/deploy-agent.service
[Unit]
Description=OpenClaw Agent Auto-Deploy + HotPool Register
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no
User=ubuntu
WorkingDirectory=/tmp

Environment=DEPLOY_SERVER=http://43.160.245.20:9900
Environment=ADMIN_API=https://www.nika8.com/api
# ↓ 部署时由 quick-deploy.sh 填入，或手动替换 ↓
Environment=ADMIN_API_KEY=
Environment=DEEPSEEK_API_KEY=
Environment=AGENT_PROVIDER=Tencent

ExecStart=/bin/bash -c ' \
  curl -sL ${DEPLOY_SERVER}/register-agent.py -o /tmp/register-agent.py && \
  curl -sL ${DEPLOY_SERVER}/deploy-and-register.sh -o /tmp/deploy-and-register.sh && \
  bash /tmp/deploy-and-register.sh'

StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
